# scripts/logging/infrastructure-logs-daemon.service
[Unit]
Description=Infrastructure Logs Daemon (Tail /var/log/* to Cassandra)
Documentation=file:///opt/iact/scripts/logging/README.md
After=network.target cassandra.service
Wants=cassandra.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/iact/scripts/logging

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Command
ExecStart=/usr/bin/python3 /opt/iact/scripts/logging/infrastructure_logs_daemon.py \
    --cassandra-hosts cassandra-1.internal cassandra-2.internal cassandra-3.internal \
    --keyspace logging \
    --batch-size 1000 \
    --batch-timeout 5.0 \
    --queue-maxsize 100000 \
    --health-port 9090

# Restart policy
Restart=always
RestartSec=10s
StartLimitInterval=0

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security (adjust as needed)
# NOTE: Requiere root para leer /var/log/*, considerar log forwarding alternativo
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=infrastructure-logs-daemon

[Install]
WantedBy=multi-user.target
